# systemd/DNA_Homescreen-launchpad-ready.path
# Description=Path activation for DNA_Homescreen

[Path]
PathExists=%t/wayland-0
Unit=DNA_Homescreen.service

[Install]
WantedBy=tizen-user-middleware.target
