# cameras_ext/jlr_cmss/service/jlr_cmss.service
[Unit]
Description=jlr_cmsss for JLRCameras app

[Service]
Type=dbus
BusName=com.jlr.JLRCameras
ExecStart=/usr/sbin/jlr_cmss
SmackExecLabel=User
